# systemd/ai-orchestrator.service
[Unit]
Description=AI Orchestrator — Central AI Coordination Daemon
Documentation=https://github.com/Khazar-System-Distribution/ai-orchestrator
After=network.target local-fs.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ai-orchestrator /etc/ai-orchestrator.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StartLimitInterval=30
StartLimitBurst=3
LimitNOFILE=4096
MemoryMax=100M
CPUQuota=50%
Nice=-5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-orchestrator

[Install]
WantedBy=multi-user.target
