# systemd/ai-orchestrator.socket
[Unit]
Description=AI Orchestrator Socket
PartOf=ai-orchestrator.service

[Socket]
ListenStream=/run/ai-orchestrator.sock
SocketMode=0660
SocketUser=root
SocketGroup=ai
RemoveOnStop=true

[Install]
WantedBy=sockets.target
