# systemd/ai-orchestrator-agent@.service
[Unit]
Description=AI Orchestrator Agent (%i)
PartOf=ai-orchestrator.service
BindsTo=ai-orchestrator.service
After=ai-orchestrator.service
Wants=ai-orchestrator.service

[Service]
Type=simple
ExecStart=/usr/local/lib/ai-orchestrator/agents/%I /run/ai-orchestrator.sock
Restart=always
RestartSec=3
StartLimitInterval=30
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-agent-%i

[Install]
WantedBy=multi-user.target
